# config/systemd/user/change-background.timer
[Unit]
Description=Change background image XX:10 each hour

[Timer]
OnCalendar=*-*-* *:10:00
Persistent=true

[Install]
WantedBy=timers.target
