# config/systemd/user/flatpak-updates.timer
[Unit]
Description=Run flatpak updates check XX:45 each hour

[Timer]
OnCalendar=*-*-* *:45:00
Persistent=true

[Install]
WantedBy=timers.target
