# config/systemd/user/fwupd-refresh.service
[Unit]
Description=Firmware Update (LVFS) synchronization service

[Service]
Type=oneshot
ExecStart=/home/mcritchlow/.local/bin/fwupd-refresh
